# daemontools/etc_init_daemontools.conf
# djb daemontools  init 
#
# This service maintains a svscanboot on the specified device.

start on runlevel [345]
stop on runlevel [S016]

respawn
exec _PREFIX_/svscanboot

